# eorlive_prototype/crontab
*/10 * * * * /mnt/eorlive_prototype/pull_histogram_data.sh
